# var/lib/dpkg/info/kmod.list
/.
/usr
/usr/share
/usr/share/man/man5/depmod.d.5.gz
/usr/share/man/man5/modprobe.d.5.gz
/usr/share/man/man5/modules.dep.5.gz
/usr/share/man/man5/modules.5.gz
/usr/share/man/man8/depmod.8.gz
/usr/share/man/man8/insmod.8.gz
/usr/share/man/man8/lsmod.8.gz
/usr/share/man/man8/rmmod.8.gz
/usr/share/man/man8/modprobe.8.gz
/usr/share/man/man8/modinfo.8.gz
/usr/share/man/fr/man5/modules.5.gz
/usr/share/doc
/usr/share/doc/libkmod2
/usr/share/doc/libkmod2/README
/usr/share/doc/libkmod2/TODO
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/kmod
/bin
/bin/kmod
/etc
/etc/modprobe.d
/etc/init.d
/etc/init.d/kmod
/lib
/lib/modprobe.d
/lib/modprobe.d/aliases.conf
/sbin
/usr/share/man/man5/modules.dep.bin.5.gz
/usr/share/doc/kmod
/bin/lsmod
/sbin/rmmod
/sbin/modprobe
/sbin/modinfo
/sbin/lsmod
/sbin/insmod
/sbin/depmod
