# var/lib/dpkg/info/util-linux.list
/.
/sbin
/sbin/wipefs
/sbin/ctrlaltdel
/sbin/swaplabel
/sbin/mkfs.cramfs
/sbin/fsck
/sbin/mkfs.minix
/sbin/fdisk
/sbin/fsck.minix
/sbin/mkswap
/sbin/agetty
/sbin/cfdisk
/sbin/switch_root
/sbin/mkfs
/sbin/sfdisk
/sbin/findfs
/sbin/fsck.cramfs
/sbin/blkid
/sbin/raw
/sbin/fsfreeze
/sbin/blockdev
/sbin/pivot_root
/sbin/fstrim
/sbin/mkfs.bfs
/sbin/isosize
/sbin/hwclock
/sbin/getty
/bin
/bin/tailf
/bin/more
/bin/lsblk
/bin/dmesg
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/85-hwclock.rules
/lib/udev/hwclock-set
/lib/x86_64-linux-gnu
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ipcmk.1.gz
/usr/share/man/man1/dmesg.1.gz
/usr/share/man/man1/rev.1.gz
/usr/share/man/man1/more.1.gz
/usr/share/man/man1/ionice.1.gz
/usr/share/man/man1/getopt.1.gz
/usr/share/man/man1/ddate.1.gz
/usr/share/man/man1/line.1.gz
/usr/share/man/man1/mcookie.1.gz
/usr/share/man/man1/setterm.1.gz
/usr/share/man/man1/linux64.1.gz
/usr/share/man/man1/chkdupexe.1.gz
/usr/share/man/man1/namei.1.gz
/usr/share/man/man1/taskset.1.gz
/usr/share/man/man1/unshare.1.gz
/usr/share/man/man1/fallocate.1.gz
/usr/share/man/man1/readprofile.1.gz
/usr/share/man/man1/ipcs.1.gz
/usr/share/man/man1/tailf.1.gz
/usr/share/man/man1/chrt.1.gz
/usr/share/man/man1/whereis.1.gz
/usr/share/man/man1/lscpu.1.gz
/usr/share/man/man1/ipcrm.1.gz
/usr/share/man/man1/rename.ul.1.gz
/usr/share/man/man1/linux32.1.gz
/usr/share/man/man1/pg.1.gz
/usr/share/man/man1/flock.1.gz
/usr/share/man/man1/setsid.1.gz
/usr/share/man/man3
/usr/share/man/ru
/usr/share/man/ru/man1
/usr/share/man/ru/man1/ddate.1.gz
/usr/share/man/man8
/usr/share/man/man8/i386.8.gz
/usr/share/man/man8/blkid.8.gz
/usr/share/man/man8/mkfs.bfs.8.gz
/usr/share/man/man8/cytune.8.gz
/usr/share/man/man8/ldattach.8.gz
/usr/share/man/man8/tunelp.8.gz
/usr/share/man/man8/rtcwake.8.gz
/usr/share/man/man8/fsfreeze.8.gz
/usr/share/man/man8/raw.8.gz
/usr/share/man/man8/partx.8.gz
/usr/share/man/man8/switch_root.8.gz
/usr/share/man/man8/fsck.minix.8.gz
/usr/share/man/man8/fdisk.8.gz
/usr/share/man/man8/agetty.8.gz
/usr/share/man/man8/mkfs.8.gz
/usr/share/man/man8/mkfs.minix.8.gz
/usr/share/man/man8/cfdisk.8.gz
/usr/share/man/man8/addpart.8.gz
/usr/share/man/man8/ctrlaltdel.8.gz
/usr/share/man/man8/sfdisk.8.gz
/usr/share/man/man8/swaplabel.8.gz
/usr/share/man/man8/lsblk.8.gz
/usr/share/man/man8/setarch.8.gz
/usr/share/man/man8/getty.8.gz
/usr/share/man/man8/fsck.8.gz
/usr/share/man/man8/pivot_root.8.gz
/usr/share/man/man8/delpart.8.gz
/usr/share/man/man8/x86_64.8.gz
/usr/share/man/man8/fstrim.8.gz
/usr/share/man/man8/hwclock.8.gz
/usr/share/man/man8/isosize.8.gz
/usr/share/man/man8/blockdev.8.gz
/usr/share/man/man8/wipefs.8.gz
/usr/share/man/man8/fdformat.8.gz
/usr/share/man/man8/findfs.8.gz
/usr/share/man/man8/mkswap.8.gz
/usr/share/man/man5
/usr/share/man/man5/hwclock.5.gz
/usr/share/info
/usr/share/doc
/usr/share/doc/util-linux
/usr/share/doc/util-linux/README.cfdisk
/usr/share/doc/util-linux/README.fdisk.gz
/usr/share/doc/util-linux/examples
/usr/share/doc/util-linux/examples/getopt-parse.tcsh
/usr/share/doc/util-linux/examples/getopt-parse.bash
/usr/share/doc/util-linux/examples/sfdisk.examples.gz
/usr/share/doc/util-linux/examples/fstab.example2
/usr/share/doc/util-linux/README.modems-with-agetty
/usr/share/doc/util-linux/README.Debian.hwclock
/usr/share/doc/util-linux/changelog.gz
/usr/share/doc/util-linux/AUTHORS.gz
/usr/share/doc/util-linux/copyright
/usr/share/doc/util-linux/README.poeigl.gz
/usr/share/doc/util-linux/changelog.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/util-linux
/usr/sbin
/usr/sbin/fdformat
/usr/sbin/ldattach
/usr/sbin/readprofile
/usr/sbin/tunelp
/usr/sbin/rtcwake
/usr/sbin/cytune
/usr/bin
/usr/bin/setarch
/usr/bin/ipcs
/usr/bin/taskset
/usr/bin/namei
/usr/bin/chkdupexe
/usr/bin/ipcmk
/usr/bin/ionice
/usr/bin/ipcrm
/usr/bin/addpart
/usr/bin/partx
/usr/bin/rev
/usr/bin/lscpu
/usr/bin/ddate
/usr/bin/rename.ul
/usr/bin/setsid
/usr/bin/setterm
/usr/bin/fallocate
/usr/bin/chrt
/usr/bin/unshare
/usr/bin/mcookie
/usr/bin/pg
/usr/bin/flock
/usr/bin/line
/usr/bin/delpart
/usr/bin/getopt
/usr/bin/whereis
/usr/lib
/usr/lib/mime
/usr/lib/mime/packages
/usr/lib/mime/packages/util-linux
/usr/lib/x86_64-linux-gnu
/etc
/etc/default
/etc/default/hwclock
/etc/fstab.d
/etc/init.d
/etc/init.d/hwclock.sh
/usr/share/man/man8/clock.8.gz
/usr/bin/i386
/usr/bin/x86_64
/usr/bin/linux32
/usr/bin/linux64
